# zpool-scrub@.timer
[Unit]
Description=Scrub ZFS pool monthly

[Timer]
OnCalendar=monthly
Persistent=true

[Install]
WantedBy=timers.target
